# system/files/discord-update.timer
[Unit]
Description=Auto-update Discord on boot

[Timer]
OnBootSec=5s

[Install]
WantedBy=timers.target
