# discord-update/etc/systemd/system/discord-update.service
[Unit]
Description=Auto-update Discord
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/bin/pacman -Sy discord --noconfirm
