# system/files/claude/claude-remote-control.service
[Unit]
Description=Claude Code Remote Control
After=network-online.target

[Service]
WorkingDirectory=/home/connor/claude-remote
# %H expands to the machine's hostname at unit-load time.
ExecStart=/home/connor/.local/bin/claude remote-control --name "%H"
Restart=on-failure
RestartSec=10

[Install]
WantedBy=default.target
