# hypr/.config/systemd/user/wallpaper-cycle.service
[Unit]
Description=Cycle Hyprland wallpaper

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/wallpaper-cycle.sh
